# tg_monitor.service
[Unit]
Description=Telegram Channel Monitor Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/tg_monitor
ExecStart=/home/ubuntu/tg_monitor/venv/bin/python main.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
